# rpmbuild/SOURCES/urpmd.service
[Unit]
Description=urpm daemon for intelligent cache management
Documentation=man:urpmd(8)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/urpmd --foreground
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
